# config/systemd/fks-ssl-manager.service
[Unit]
Description=FKS SSL Certificate Manager
Documentation=https://github.com/nuniesmith/fks
After=network-online.target nginx.service
Wants=network-online.target
Requires=nginx.service

[Service]
Type=oneshot
User=root
Group=root
Environment=DOMAIN=fkstrading.xyz
Environment=LETSENCRYPT_EMAIL=nunie.smith01@gmail.com
Environment=WEBROOT_PATH=/var/www/html
Environment=STAGING=false
ExecStart=/usr/local/bin/fks_ssl-manager.sh install
ExecReload=/usr/local/bin/fks_ssl-manager.sh renew
StandardOutput=journal
StandardError=journal
TimeoutSec=300

# Security settings
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=false
ProtectHome=true
PrivateDevices=false

[Install]
WantedBy=multi-user.target
